# system_files/etc/profile.d/editor.sh
export EDITOR="${EDITOR:-nano}"
export VISUAL="${VISUAL:-nano}"
